# systemd/dotclaw.service
[Unit]
Description=DotClaw Telegram Assistant
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{DOTCLAW_HOME}}
Environment=NODE_ENV=production
Environment=DOTCLAW_HOME={{DOTCLAW_HOME}}
EnvironmentFile={{DOTCLAW_HOME}}/.env
ExecStart={{NODE_PATH}} {{PROJECT_ROOT}}/dist/index.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillMode=process
StandardOutput=append:{{DOTCLAW_HOME}}/logs/dotclaw.log
StandardError=append:{{DOTCLAW_HOME}}/logs/dotclaw.error.log

[Install]
WantedBy=multi-user.target
